# 99-cd-audio-processing.rules
SUBSYSTEM=="block", KERNEL=="sr0", ENV{ID_CDROM_MEDIA_CD}=="1", ACTION=="change", RUN+="/bin/systemctl start ripit.service"
